# NJUNLP-SMT/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:tornado_NJU]
command=python /home/longsail/NJUNLP-WebDemo/NJUNLP-SMT/main.py --port=1989
directory = /home/longsail/NJUNLP-WebDemo/NJUNLP-SMT
autorestart=true
redirect_stderr=true
user=www-data
stdout_logfile=/home/longsail/NJUNLP-WebDemo/NJUNLP-SMT/logs/tornado.log


